# experiments/xcapture-v0-proc/xcapture-restart.service
[Unit]
Description=restart xcapture

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart xcapture
