# xcapture-restart.timer
[Unit]
Description=xcapture restart

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
