# xcapture/packaging/systemd/xcapture.service.in
[Unit]
Description=xcapture system activity monitor
After=network.target
Documentation=https://github.com/0xtools/0xtools-next

[Service]
Type=simple
EnvironmentFile=-/etc/xcapture.conf
User=xcapture
Group=xcapture
PermissionsStartOnly=true
ExecStartPre=/usr/bin/install -d -m 0750 -o xcapture -g xcapture /var/log/xcapture
ExecStartPre=/usr/bin/install -d -m 0750 -o xcapture -g xcapture /var/lib/xcapture
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/xcapture $XCAPTURE_OPTS
WorkingDirectory=/var/lib/xcapture
Restart=on-failure
RestartSec=5s
StandardOutput=append:/var/log/xcapture/xcapture.log
StandardError=append:/var/log/xcapture/xcapture.log
AmbientCapabilities=CAP_BPF CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_BPF CAP_SYS_ADMIN
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
PrivateDevices=yes
PrivateTmp=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
ReadWritePaths=/var/log/xcapture /var/lib/xcapture /sys/fs/bpf

[Install]
WantedBy=multi-user.target
